# roles/sonm-geth/templates/geth.service.j2
[Unit]
Description=Geth node
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=30s
ExecStart=/usr/bin/geth --datadir {{ geth.geth_data_dir }} --networkid 444 --cache 512  --maxpeers 2048 --rpc --rpcaddr localhost --rpcport 8545 --rpcapi "web3,eth,net" --rpccorsdomain "http://localhost:8000" --bootnodes="enode://ad8c2b6767285f6cf0678bee20d304a85a1336684e606c4c87ab3941782aedfdc67eec912c93d33b72b1d5192d2301d3ebb10aff7abe8e6fe518d627b639497f@159.65.167.139:30303,enode://ad8c2b6767285f6cf0678bee20d304a85a1336684e606c4c87ab3941782aedfdc67eec912c93d33b72b1d5192d2301d3ebb10aff7abe8e6fe518d627b639497f@159.65.167.139:30313,enode://a8b24454c7da48066ca6e72d3dcc986775318461a417bb17b0c89e1ba272b21b3b770fca0a5faf449d324e64ab0110b1352ed2355cb908bf506439dac9ce7d5d@159.89.205.255:30303,enode://a8b24454c7da48066ca6e72d3dcc986775318461a417bb17b0c89e1ba272b21b3b770fca0a5faf449d324e64ab0110b1352ed2355cb908bf506439dac9ce7d5d@159.89.205.255:30313,enode://8e83114a70453e0a8e556744da07c50724bd40cb846fea8fe4ce3d8129ab51e5ac8208468fc30e1e8a143e568994ae4c70b98bb9e4b58a48e52d4c7085c4bb21@46.101.241.114:30303,enode://8e83114a70453e0a8e556744da07c50724bd40cb846fea8fe4ce3d8129ab51e5ac8208468fc30e1e8a143e568994ae4c70b98bb9e4b58a48e52d4c7085c4bb21@46.101.241.114:30313,enode://1ca0323bd62f3489f1fb00091fe4a85c1d74931a25bcd7a99232beb60a809e2e5a3dde461a68bcb73202d0c7b0c7fe8f0df82342c35a221b08e89099763e4560@159.65.22.22:30303,enode://1ca0323bd62f3489f1fb00091fe4a85c1d74931a25bcd7a99232beb60a809e2e5a3dde461a68bcb73202d0c7b0c7fe8f0df82342c35a221b08e89099763e4560@159.65.22.22:30313"

[Install]
WantedBy=default.target
